# tomcat/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/tomcat
drwxr-xr-x	root/root	etc/tomcat/
-rw-------	root/root	etc/tomcat/catalina.policy
-rw-------	root/root	etc/tomcat/catalina.properties
-rw-------	root/root	etc/tomcat/context.xml
-rw-------	root/root	etc/tomcat/logging.properties
-rw-------	root/root	etc/tomcat/server.xml
-rw-------	root/root	etc/tomcat/tomcat-users.xml
-rw-------	root/root	etc/tomcat/web.xml
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/libtcnative-1.a
-rwxr-xr-x	root/root	usr/lib/libtcnative-1.la
lrwxrwxrwx	root/root	usr/lib/libtcnative-1.so -> libtcnative-1.so.0.1.16
lrwxrwxrwx	root/root	usr/lib/libtcnative-1.so.0 -> libtcnative-1.so.0.1.16
-rwxr-xr-x	root/root	usr/lib/libtcnative-1.so.0.1.16
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/tcnative-1.pc
drwxr-xr-x	root/root	usr/lib/tomcat/
-rw-r--r--	root/root	usr/lib/tomcat/annotations-api.jar
-rw-r--r--	root/root	usr/lib/tomcat/catalina-ant.jar
-rw-r--r--	root/root	usr/lib/tomcat/catalina-ha.jar
-rw-r--r--	root/root	usr/lib/tomcat/catalina-tribes.jar
-rw-r--r--	root/root	usr/lib/tomcat/catalina.jar
-rw-r--r--	root/root	usr/lib/tomcat/el-api.jar
-rw-r--r--	root/root	usr/lib/tomcat/jasper-el.jar
-rw-r--r--	root/root	usr/lib/tomcat/jasper-jdt.jar
-rw-r--r--	root/root	usr/lib/tomcat/jasper.jar
-rw-r--r--	root/root	usr/lib/tomcat/jsp-api.jar
-rw-r--r--	root/root	usr/lib/tomcat/servlet-api.jar
-rw-r--r--	root/root	usr/lib/tomcat/tomcat-coyote.jar
-rw-r--r--	root/root	usr/lib/tomcat/tomcat-dbcp.jar
drwxr-xr-x	root/root	usr/opt/
drwxr-xr-x	root/root	usr/opt/tomcat/
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/jsvc
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/tomcat/
drwxr-xr-x	root/root	usr/share/tomcat/bin/
-rw-r--r--	root/root	usr/share/tomcat/bin/bootstrap.jar
-rw-r--r--	root/root	usr/share/tomcat/bin/catalina-tasks.xml
-rwxr-xr-x	root/root	usr/share/tomcat/bin/catalina.sh
-rw-r--r--	root/root	usr/share/tomcat/bin/commons-daemon.jar
-rwxr-xr-x	root/root	usr/share/tomcat/bin/digest.sh
-rwxr-xr-x	root/root	usr/share/tomcat/bin/setclasspath.sh
-rwxr-xr-x	root/root	usr/share/tomcat/bin/shutdown.sh
-rwxr-xr-x	root/root	usr/share/tomcat/bin/startup.sh
-rw-r--r--	root/root	usr/share/tomcat/bin/tomcat-juli.jar
-rwxr-xr-x	root/root	usr/share/tomcat/bin/tool-wrapper.sh
-rwxr-xr-x	root/root	usr/share/tomcat/bin/version.sh
lrwxrwxrwx	root/root	usr/share/tomcat/conf -> ../../../etc/tomcat
lrwxrwxrwx	root/root	usr/share/tomcat/lib -> ../../../usr/lib/tomcat
lrwxrwxrwx	root/root	usr/share/tomcat/logs -> ../../../var/log/tomcat
lrwxrwxrwx	root/root	usr/share/tomcat/temp -> ../../../var/run/tomcat/temp
drwxr-xr-x	root/root	usr/share/tomcat/webapps/
drwxr-xr-x	root/root	usr/share/tomcat/webapps/ROOT/
-rw-r--r--	root/root	usr/share/tomcat/webapps/ROOT/RELEASE-NOTES.txt
drwxr-xr-x	root/root	usr/share/tomcat/webapps/ROOT/WEB-INF/
-rw-r--r--	root/root	usr/share/tomcat/webapps/ROOT/WEB-INF/web.xml
-rw-r--r--	root/root	usr/share/tomcat/webapps/ROOT/asf-logo-wide.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/ROOT/build.xml
-rw-r--r--	root/root	usr/share/tomcat/webapps/ROOT/favicon.ico
-rw-r--r--	root/root	usr/share/tomcat/webapps/ROOT/index.html
-rw-r--r--	root/root	usr/share/tomcat/webapps/ROOT/index.jsp
-rw-r--r--	root/root	usr/share/tomcat/webapps/ROOT/tomcat-power.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/ROOT/tomcat.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/ROOT/tomcat.svg
drwxr-xr-x	root/root	usr/share/tomcat/webapps/host-manager/
drwxr-xr-x	root/root	usr/share/tomcat/webapps/host-manager/META-INF/
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/META-INF/context.xml
drwxr-xr-x	root/root	usr/share/tomcat/webapps/host-manager/WEB-INF/
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/WEB-INF/web.xml
drwxr-xr-x	root/root	usr/share/tomcat/webapps/host-manager/images/
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/images/add.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/images/asf-logo.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/images/code.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/images/design.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/images/docs.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/images/fix.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/images/tomcat.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/images/update.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/images/void.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/host-manager/manager.xml
drwxr-xr-x	root/root	usr/share/tomcat/webapps/manager/
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/401.jsp
drwxr-xr-x	root/root	usr/share/tomcat/webapps/manager/META-INF/
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/META-INF/context.xml
drwxr-xr-x	root/root	usr/share/tomcat/webapps/manager/WEB-INF/
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/WEB-INF/web.xml
drwxr-xr-x	root/root	usr/share/tomcat/webapps/manager/images/
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/images/add.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/images/asf-logo.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/images/code.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/images/design.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/images/docs.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/images/fix.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/images/tomcat.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/images/update.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/images/void.gif
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/sessionDetail.jsp
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/sessionsList.jsp
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/status.xsd
-rw-r--r--	root/root	usr/share/tomcat/webapps/manager/xform.xsl
lrwxrwxrwx	root/root	usr/share/tomcat/work -> ../../../var/run/tomcat/work
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/log/tomcat/
drwxr-xr-x	root/root	var/run/
drwxr-xr-x	root/root	var/run/tomcat/
drwxr-xr-x	root/root	var/run/tomcat/temp/
drwxr-xr-x	root/root	var/run/tomcat/work/
