# postgrey/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/postfix/
-rw-r--r--	root/root	etc/postfix/postgrey_whitelist_clients
-rw-r--r--	root/root	etc/postfix/postgrey_whitelist_recipients
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/postgrey
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/postgrey
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/spool/
drwxr-xr-x	root/root	var/spool/postfix/
drwxr-xr-x	postgrey/postgrey	var/spool/postfix/postgrey/
