# icecast/.footprint
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/icecast.xml
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/icecast
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/icecast
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/icecast/
drwxr-xr-x	root/root	usr/share/icecast/admin/
-rw-r--r--	root/root	usr/share/icecast/admin/listclients.xsl
-rw-r--r--	root/root	usr/share/icecast/admin/listmounts.xsl
-rw-r--r--	root/root	usr/share/icecast/admin/manageauth.xsl
-rw-r--r--	root/root	usr/share/icecast/admin/moveclients.xsl
-rw-r--r--	root/root	usr/share/icecast/admin/response.xsl
-rw-r--r--	root/root	usr/share/icecast/admin/stats.xsl
-rw-r--r--	root/root	usr/share/icecast/admin/updatemetadata.xsl
drwxr-xr-x	root/root	usr/share/icecast/web/
-rw-r--r--	root/root	usr/share/icecast/web/auth.xsl
-rw-r--r--	root/root	usr/share/icecast/web/corner_bottomleft.jpg
-rw-r--r--	root/root	usr/share/icecast/web/corner_bottomright.jpg
-rw-r--r--	root/root	usr/share/icecast/web/corner_topleft.jpg
-rw-r--r--	root/root	usr/share/icecast/web/corner_topright.jpg
-rw-r--r--	root/root	usr/share/icecast/web/icecast.png
-rw-r--r--	root/root	usr/share/icecast/web/key.gif
-rw-r--r--	root/root	usr/share/icecast/web/key.png
-rw-r--r--	root/root	usr/share/icecast/web/server_version.xsl
-rw-r--r--	root/root	usr/share/icecast/web/status.xsl
-rw-r--r--	root/root	usr/share/icecast/web/status2.xsl
-rw-r--r--	root/root	usr/share/icecast/web/style.css
-rw-r--r--	root/root	usr/share/icecast/web/tunein.png
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/log/icecast/
-rw-r--r--	root/root	var/log/icecast/access.log (EMPTY)
-rw-r--r--	root/root	var/log/icecast/error.log (EMPTY)
